# scripts/systemd/polyhermes-bridge.service
[Unit]
Description=PolyHermes Polymtrade Bridge
After=network.target

[Service]
Type=simple
User=polyhermes
Group=polyhermes
WorkingDirectory=/opt/polyhermes/polymtrade-bridge
Environment="BRIDGE_PORT=8080"
EnvironmentFile=-/etc/polyhermes/bridge.env
ExecStart=/opt/polyhermes/polymtrade-bridge/start.sh
ExecStop=/bin/kill -SIGTERM $MAINPID
TimeoutStopSec=60
Restart=on-failure
RestartSec=30

[Install]
WantedBy=multi-user.target
