# scripts/launchd/com.polyhermes.bridge.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.polyhermes.bridge</string>

    <key>ProgramArguments</key>
    <array>
        <string>/bin/bash</string>
        <string>/Users/henry/projects/polyhermes/polymtrade-bridge/start.sh</string>
    </array>

    <key>EnvironmentVariables</key>
    <dict>
        <key>BRIDGE_PORT</key>
        <string>8080</string>
    </dict>

    <key>WorkingDirectory</key>
    <string>/Users/henry/projects/polyhermes/polymtrade-bridge</string>

    <key>StandardOutPath</key>
    <string>/Users/henry/projects/polyhermes/logs/bridge.log</string>
    <key>StandardErrorPath</key>
    <string>/Users/henry/projects/polyhermes/logs/bridge.log</string>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
        <key>Crashed</key>
        <true/>
    </dict>

    <key>ThrottleInterval</key>
    <integer>30</integer>

    <key>RunAtLoad</key>
    <true/>
</dict>
</plist>
